# packaging/systemd/yadiskd.service
[Unit]
Description=Yandex Disk sync daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/yadiskd
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
Environment=YADISK_SYNC_DIR=%h/Yandex Disk
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yadiskd

[Install]
WantedBy=default.target
